# testing/integration/test/etc/supervisord-firefox.conf
[supervisord]

[program:firefox]
directory=/test/src/uproxy/build/dev/uproxy/lib/samples/zork-firefoxapp
command=jpm run -b /usr/bin/firefox
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/zork.log
