# testing/integration/test/etc/supervisord-chrome.conf
[supervisord]

[program:chrome]
command=google-chrome --user-data-dir=/tmp/chrome-data --load-and-launch-app=/test/src/uproxy-lib/build/dev/uproxy-lib/samples/zork-chromeapp --no-default-browser-check --no-first-run
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/zork.log
