# testing/integration/test/etc/supervisord-chrome.conf
[supervisord]
# Without this, Chrome frequently hangs on startup.
# We use a workaround suggested in this Selenium issue:
#   https://github.com/SeleniumHQ/docker-selenium/issues/87#issuecomment-187580115
environment=DBUS_SESSION_BUS_ADDRESS="/dev/null"

[program:chrome]
command=google-chrome --user-data-dir=/tmp/chrome-data --load-and-launch-app=/test/zork/zork-chromeapp --no-default-browser-check --no-first-run
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/zork.log
