# docker/sshd/sshd.ini
[program:sshd]
command=bash -c "ssh-keygen -A && /usr/sbin/sshd -D -e"
stderr_logfile=/var/log/sshd.log
