# docker/sshd/supervisord.conf
[program:sshd]
command=/usr/sbin/sshd -D