# docker/testing/integration/test/etc/supervisord-node.conf
[supervisord]

[program:node]
command=nodejs /test/zork/zork-node
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/zork.log
