# docker/testing/integration/test/etc/supervisord-firefox.conf
[supervisord]

[program:firefox]
directory=/test/zork/zork-firefoxapp
command=jpm run -b /usr/bin/firefox
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/zork.log
